# supervisord.conf
[supervisord]
nodaemon=true

[program:couchpotato]
user=couchpotato
command=python /opt/couchpotato/CouchPotato.py --data_dir=/config
